# dimmer.service
[Unit]
Description=A simple backlight dimmer

[Service]
ExecStart=/usr/bin/dimmer

[Install]
WantedBy=multi-user.target
